# buttcoin-master/contrib/init/bitcoind.conf
description "buttcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env buttcoinD_BIN="/usr/bin/buttcoind"
env buttcoinD_USER="buttcoin"
env buttcoinD_GROUP="buttcoin"
env buttcoinD_PIDDIR="/var/run/buttcoind"
# upstart can't handle variables constructed with other variables
env buttcoinD_PIDFILE="/var/run/buttcoind/buttcoind.pid"
env buttcoinD_CONFIGFILE="/etc/buttcoin/buttcoin.conf"
env buttcoinD_DATADIR="/var/lib/buttcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # buttcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$buttcoinD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run buttcoind."
        echo "The setting must appear in $buttcoinD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"buttcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$buttcoinD_PIDDIR"
    chmod 0755 "$buttcoinD_PIDDIR"
    chown $buttcoinD_USER:$buttcoinD_GROUP "$buttcoinD_PIDDIR"
    chown $buttcoinD_USER:$buttcoinD_GROUP "$buttcoinD_CONFIGFILE"
    chmod 0660 "$buttcoinD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$buttcoinD_PIDFILE" \
    --chuid $buttcoinD_USER:$buttcoinD_GROUP \
    --exec "$buttcoinD_BIN" \
    -- \
    -pid="$buttcoinD_PIDFILE" \
    -conf="$buttcoinD_CONFIGFILE" \
    -datadir="$buttcoinD_DATADIR" \
    -disablewallet \
    -daemon

